# configurations/gnome-autostart/alarm-clock.desktop
[Desktop Entry]
Type=Application
Encoding=UTF-8
Version=1.0
Name=Alarm Clock
Exec=alarm-clock
X-GNOME-Autostart-enabled=true
